# service/log/run
#!/bin/bash
exec multilog t n20 s1000000 ./main
